# deploy/weibo-flask.service
[Unit]
Description=Weibo Archive Web Server - Production WSGI server
After=network.target

[Service]
Type=notify
User=%i
WorkingDirectory=/home/%i/weibo-archive
EnvironmentFile=/home/%i/weibo-archive/.env

# 使用Gunicorn运行Flask（生产级WSGI服务器）
# -w 4: 4个worker进程（根据CPU核心数调整）
# -b 127.0.0.1:5000: 仅监听本地，通过Nginx反向代理
# --timeout 120: 请求超时时间（秒）
# --access-logfile: 访问日志
# --error-logfile: 错误日志
ExecStart=/home/%i/weibo-archive/venv/bin/gunicorn \
    --workers 4 \
    --bind 127.0.0.1:5000 \
    --timeout 120 \
    --access-logfile /home/%i/weibo-archive/logs/gunicorn_access.log \
    --error-logfile /home/%i/weibo-archive/logs/gunicorn_error.log \
    --log-level info \
    app:app

# 自动重启
Restart=always
RestartSec=5

# 进程管理
KillMode=mixed
KillSignal=SIGQUIT

# 日志
StandardOutput=append:/home/%i/weibo-archive/logs/flask.log
StandardError=append:/home/%i/weibo-archive/logs/flask_error.log

[Install]
WantedBy=multi-user.target
