# deploy/weibo-scheduler@.service
[Unit]
Description=Weibo Archive Scheduler - Smart crawling service
After=network.target

[Service]
Type=simple
User=%i
WorkingDirectory=/home/%i/WeiboCrawler
EnvironmentFile=/home/%i/WeiboCrawler/.env

# 激活虚拟环境并运行调度器
ExecStart=/home/%i/WeiboCrawler/venv/bin/python /home/%i/WeiboCrawler/scheduler.py

# 自动重启
Restart=always
RestartSec=10

# 日志
StandardOutput=append:/home/%i/WeiboCrawler/logs/scheduler.log
StandardError=append:/home/%i/WeiboCrawler/logs/scheduler_error.log

[Install]
WantedBy=multi-user.target
